# rootfs/etc/inittab
::sysinit:/etc/init.d/rcS
#::respawn:-/bin/sh
::respawn:/sbin/getty 115200 ttyO0
::ctrlaltdel:/bin/umount -a -r
